# docs/systemd-scripts/snb-enrollment-wecom-worker.service
[Unit]
Description=SNB Enrollment WeCom Worker
After=network.target

[Service]
Type=simple
User=administrator
Group=administrator
WorkingDirectory=/home/administrator/Desktop/snb-enrollment/backend
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=-/etc/snb-enrollment/common.env
EnvironmentFile=-/etc/snb-enrollment/wecom-worker.env
ExecStart=/home/administrator/Desktop/snb-enrollment/.venv/bin/python -m app.workers.wecom_message_worker
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
